# config/ldap-matching-report.service
[Unit]
Description=LDAP-Locosoft Matching Report (wöchentlich)

[Service]
Type=oneshot
User=ag-admin
WorkingDirectory=/opt/greiner-portal
ExecStart=/opt/greiner-portal/venv/bin/python3 /opt/greiner-portal/scripts/ldap_locosoft_matching_report.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
