# config/systemd/celery-beat.service
[Unit]
Description=GREINER DRIVE Celery Beat Scheduler
After=network.target redis.service

[Service]
Type=simple
User=ag-admin
Group=ag-admin
WorkingDirectory=/opt/greiner-portal
Environment="PATH=/opt/greiner-portal/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/greiner-portal"
ExecStart=/opt/greiner-portal/venv/bin/celery -A celery_app beat --loglevel=info -S redbeat.RedBeatScheduler
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-beat

[Install]
WantedBy=multi-user.target
